# files/etc/uci-defaults/99_my_defaults
uci set network.wan.proto='pppoe'
uci set network.wan.username='07349947185'
uci set network.wan.password='19976745819'
uci set luci.themes.ThemeName='argon'
uci commit
