# files/etc/uci-defaults/99_my_defaults
uci set network.wan.proto='pppoe'
uci set network.wan.username='$ACCOUNT'
uci set network.wan.password='$PASSWORD'
uci set luci.themes.ThemeName='argon'
uci commit
